# services/omni-lidar.service
[Unit]
Description=Omni Robot RPLIDAR C1
After=network.target
Wants=network.target

[Service]
Type=simple
User=noob
EnvironmentFile=-/etc/omni-robot/omni.env
ExecStart=/bin/bash -lc 'exec "${OMNI_HOME:-/home/noob/omni-pi}/tools/run_lidar.sh"'
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
