# services/omni-mux.service
[Unit]
Description=Omni Robot Command Mux
After=network.target omni-can.service omni-bridge.service
Wants=network.target omni-can.service omni-bridge.service

[Service]
Type=simple
User=noob
EnvironmentFile=-/etc/omni-robot/omni.env
ExecStart=/bin/bash -lc 'exec "${OMNI_HOME:-/home/noob/omni-pi}/tools/run_cmd_mux.sh"'
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
